# deploy/celery-worker.service
[Unit]
# Fichier de référence — même logique mais autre nom de service sur le système réel :
# après setup.sh : systemctl utilise « ipxe-celery », pas celery-worker.service.
Description=iPXE Manager — Celery worker (ISO extraction / menu regeneration)
After=network.target redis-server.service ipxe-manager.service
Wants=redis-server.service

[Service]
Type=exec
User=ipxe
Group=ipxe
WorkingDirectory=/srv/ipxe/app
EnvironmentFile=/srv/ipxe/app/.env
ExecStart=/srv/ipxe/venv/bin/celery \
    -A app.tasks.celery_app \
    worker \
    --loglevel=info \
    --concurrency=2 \
    --hostname=worker@%%h \
    --logfile=/var/log/ipxe-manager/celery.log
Restart=always
RestartSec=10

# Même politique que l’unité « ipxe-celery » installée par deploy/setup.sh.

[Install]
WantedBy=multi-user.target
