# deploy/ipxe-manager.service
[Unit]
# Fichier de référence — l’installation standard installe aussi une unité similaire
# via deploy/setup.sh (avec logs dans /var/log/ipxe-manager/).
Description=iPXE Manager — FastAPI web application
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=exec
User=ipxe
Group=ipxe
WorkingDirectory=/srv/ipxe/app
EnvironmentFile=/srv/ipxe/app/.env
ExecStart=/srv/ipxe/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info \
    --access-log
Restart=always
RestartSec=5

# Même politique que l’unité générée par deploy/setup.sh (écritures sous /srv/ipxe).

[Install]
WantedBy=multi-user.target
